# com.pm2.startup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.pm2.startup</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>source ~/.zshrc 2>/dev/null || source ~/.bash_profile 2>/dev/null; export PATH="$PATH:/opt/homebrew/bin"; cd /Users/thomasbcox/Projects/thomas-writing-assistant && pm2 resurrect || (pm2 start ecosystem.config.cjs && pm2 save)</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/tmp/pm2-startup.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/pm2-startup.error.log</string>
</dict>
</plist>
